# odoo/tshark-trivy-promtail/trivy-scan.timer
[Unit]
Description=Run Trivy Scan Every Hour

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
Unit=trivy-scan.service
Persistent=true

[Install]
WantedBy=timers.target
